# root/etc/services.d/urbackup/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid urbackup /usr/local/bin/urbackupsrv run
